# src/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
user=root
command=/usr/sbin/sshd -D
